# systemd/ytstream-stream.service
[Unit]
Description=YTStream Streaming Service
After=ytstream-audio.service default.target

[Service]
WorkingDirectory=/opt/ytstream
ExecStart=/opt/ytstream/svc_stream.sh
Restart=always

[Install]
WantedBy=default.target
